# config/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/usr/bin/mysqld_with_init

[program:configure-cron]
command=/usr/local/bin/configure-cron
